# deploy/olt-streamlit.service
[Unit]
Description=OLT NOC Dashboard (Streamlit)
After=network.target

[Service]
WorkingDirectory=/home/pisujo/olt-web-ui-huawei
# STREAMLIT_PORT y STREAMLIT_BASE_PATH vienen del .env (EnvironmentFile).
# STREAMLIT_PORT debe ser el puerto real de Streamlit (8502), NO el de nginx (8501).
ExecStart=/home/pisujo/venv/bin/python3 -m streamlit run dashboard/app.py \
  --server.port ${STREAMLIT_PORT} \
  --server.headless true \
  --server.enableCORS false \
  --server.enableXsrfProtection false \
  --browser.gatherUsageStats false \
  --server.sslCertFile /home/pisujo/olt-web-ui-huawei/fullchain.pem \
  --server.sslKeyFile /home/pisujo/olt-web-ui-huawei/privkey.pem \
  --server.baseUrlPath ${STREAMLIT_BASE_PATH}
Restart=always
RestartSec=5
Environment=PATH=/home/pisujo/venv/bin:/usr/bin
EnvironmentFile=/home/pisujo/olt-web-ui-huawei/.env
KillMode=control-group

[Install]
WantedBy=multi-user.target
